# lib/capistrano/configs/upstart.conf.erb
description "Bluepill"

start on runlevel [2]
stop on runlevel [016]

expect daemon
exec /home/ubuntu/.rvm/wrappers/default/bluepill load <%= shared_path %>/config/<%= fetch(:application)%>.pill

respawn